# hass-temzit/rootfs/etc/services.d/temzit/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare temzit_ip
declare mqtt_server
declare mqtt_user
declare mqtt_password


## Get the 'message' key from the user config options.
temzit_ip=$(bashio::config 'temzit_ip')

MQTT_HOST=$(bashio::services mqtt "host")
MQTT_USER=$(bashio::services mqtt "username")
MQTT_PASSWORD=$(bashio::services mqtt "password")

## Print the message the user supplied, defaults to "<empty>"
bashio::log.info "temzit_ip: ${temzit_ip}"
bashio::log.info "MQTT_HOST: ${MQTT_HOST}"
bashio::log.info "MQTT_USER: ${MQTT_USER}"
bashio::log.info "MQTT_PASSWORD: ${MQTT_PASSWORD}"

#/usr/bin/env
#cat /data/options.json

## Run your program
cd /usr/bin
exec /usr/bin/Temzit.Console --TemzitIp=${temzit_ip} --MqttServer=${MQTT_HOST} --MqttUser=${MQTT_USER} --MqttPass=${MQTT_PASSWORD}"
